# etc/systemd/system/nvidia-undervolt.service
[Unit]
Description=Undervolt the first available Nvidia GPU device

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/%N

[Install]
WantedBy=graphical.target
